# build/my-app.service
[Unit]
Description=Java Application

[Service]
ExecStart=/usr/bin/java -jar /opt/api-test-monolith-1.0-SNAPSHOT.jar

[Install]
WantedBy=multi-user.target