# build/my-app.service
[Unit]
Description=Java Application

[Service]
ExecStart=/usr/bin/java -jar /opt/${var.jar_file}.jar

[Install]
WantedBy=multi-user.target